# systemd/dynhost.timer
[Unit]
Description=Run dynhost on boot and hourly

[Timer]
Unit=dynhost.service
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
